# puppet/modules/profile/templates/node_exporter.service.epp
<%- | String $listen_address | -%>
# Managed by Puppet -- do not edit
[Unit]
Description=Prometheus Node Exporter
After=network-online.target
Wants=network-online.target

[Service]
User=node_exporter
Group=node_exporter
Type=simple
ExecStart=/usr/local/bin/node_exporter \
    --web.listen-address=<%= $listen_address %>
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
